# scripts/server/com.mcp-memory-service.http-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mcp-memory-service.http-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/hkr/.local/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>scripts/server/run_http_server.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/hkr/GitHub/mcp-memory-service</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/hkr/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/hkr</string>
        <key>MCP_HTTP_ENABLED</key>
        <string>true</string>
        <key>MCP_HTTPS_ENABLED</key>
        <string>true</string>
        <key>MCP_MEMORY_STORAGE_BACKEND</key>
        <string>hybrid</string>
        <key>MCP_MEMORY_SQLITE_PRAGMAS</key>
        <string>journal_mode=WAL,busy_timeout=15000,cache_size=20000</string>
        <key>HF_HOME</key>
        <string>/Users/hkr/.cache/huggingface</string>
        <key>TRANSFORMERS_CACHE</key>
        <string>/Users/hkr/.cache/huggingface</string>
        <key>MCP_GRAPH_STORAGE_MODE</key>
        <string>graph_only</string>
    </dict>

    <!-- Load env file at startup via wrapper is not natively supported;
         the script itself loads .env via python-dotenv -->

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/hkr/Library/Logs/mcp-memory-http-server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/hkr/Library/Logs/mcp-memory-http-server.error.log</string>
</dict>
</plist>
